# apps/hubstats/root/etc/cron.d/dockerhub-stats
12 * * * * root /bin/bash -l -c 'cd /opt/dockerhub-stats; ruby stats_pull.rb 2>&1'
