# lsio-port/processed/hubstats/root/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash

/usr/bin/influxd -pidfile /var/run/influxdb/influxd.pid -config /config/etc/influxdb/influxdb.conf