# apps/hubstats/root/etc/services.d/sinatra/run
#!/usr/bin/with-contenv bash

cd /opt/dockerhub-stats/ || exit

ruby stats.rb
